# files/pige/pige.cron.d
#
# cron-jobs for pige
#

MAILTO=root
*/15 * * * *     root   /usr/share/pige/bin/pige-cron
